# plugins/pycharm/scripts/supervisord-pycharm.conf
[program:pycharm]
command=/usr/local/bin/start-pycharm.sh
priority=10
autostart=true
autorestart=false
; startsecs=0: supervisord considers the program RUNNING immediately.
; Without this, a crash before 10 s is classified as a startup failure,
; which triggers startretries regardless of autorestart=false.
startsecs=0
; startretries=0: do not retry startup failures at all.
startretries=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=false
